# roles/ansible/templates/ans-between
#!/bin/sh
### BEGIN INIT INFO
# Provides:          skeleton
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: ans-between
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.  This example start a
#                    single forking daemon capable of writing a pid
#                    file.  To get other behavoirs, implemend
#                    do_start(), do_stop() or other functions to
#                    override the defaults in /lib/init/init-d-script.
### END INIT INFO

# Author: Philip Bove <pgbson@gmail.com>
#
# Please remove the "Author" lines above and replace them
# with your own name if you copy and modify this script.
PROG="ans-between"
DESC="ans-between"
#path to app
PROG_PATH=/home/server/ans-between/app.py
PID_PATH="/var/run/$PROG.pid"
start(){
  #PIDS are process id files to id if program is running
  if [ -e "$PID_PATH" ]; then
    echo "$PROG is already running!" #1>$2
    exit 1
  else
    $PROG_PATH &
    echo "$PROG started"
    touch "$PID_PATH"
  fi
}
stop(){
  #echo "$PID_PATH/$PROG.pid"
  if [ -e "$PID_PATH" ]; then
    pkill -f  "gunicorn"
    rm "$PID_PATH"
    echo "PROG stopped"
  else
    echo "$PROG not running"
    exit 1
  fi
}

if [ "$(id -u)" != "0" ]; then
  echo "This script must be run as root" #1>$2
  exit 1
fi

case "$1" in
  start)
    start
    exit 0
  ;;
  stop)
    stop
    exit 0
  ;;
  restart)
    stop
    start
    exit 0
  ;;
  **)
      echo  "Usage: $0 {start|stop|restart}" #1>$2
      exit 1
  ;;
esac
